# platform/broadcom/sonic-platform-modules-ruijie/debian/platform-modules-ruijie-b6580-48cq8qc.postinst
#!/bin/sh
# postinst

kernel_version=$(uname -r)

if [ -e /boot/System.map-${kernel_version} ]; then
    depmod -a -F /boot/System.map-${kernel_version} ${kernel_version} || true
fi

#DEBHELPER#
# Enable Ruijie-b6580-48cq8qc-platform-service
depmod -a 
systemctl enable platform-modules-b6580-48cq8qc.service
systemctl start platform-modules-b6580-48cq8qc.service

#DEBHELPER#