# platform/broadcom/sonic-platform-modules-ruijie/common/self_detect/common_factest/service/sdkmsg.service
[Unit]
Description=baresdk plus socket
After=opennsl-modules.service

[Service]
ExecStart=/usr/local/bin/sdkmsg.sh
KillMode=process
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
